# bo24.conf
[program:bo24]
command=/home/ubuntu/bo24/myenv/bin/python /home/ubuntu/bo24/bot.py
directory=/home/ubuntu/bo24/
autostart=true
autorestart=true
startretries=5
numprocs=1
startsecs=0
stderr_logfile=/var/log/supervisor/%(program_name)s_stderr.log
stderr_logfile_maxbytes=10MB
stdout_logfile=/var/log/supervisor/%(program_name)s_stdout.log
stdout_logfile_maxbytes=10MB
environment=PATH="/home/ubuntu/bo24/myenv/bin:%(ENV_PATH)s"
